# .config/systemd/user/backup_borg.timer
[Unit]
Description=Automatic borg backup

[Timer]
OnBootSec=6h
OnUnitInactiveSec=6h
Unit=backup_borg.service

[Install]
WantedBy=timers.target
